# service/com.productionhub.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Production Hub Watchdog — launchd user agent

  Starts the watchdog at login. The watchdog in turn manages the
  hub process (starts it, health-checks it, restarts if wedged).

  launchd also handles:
    - KeepAlive: restarts the *watchdog* if it crashes
    - ThrottleInterval: prevents restart loops (min 30s between restarts)

  Install:
    cp service/com.productionhub.watchdog.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.productionhub.watchdog.plist

  Uninstall:
    launchctl unload ~/Library/LaunchAgents/com.productionhub.watchdog.plist
    rm ~/Library/LaunchAgents/com.productionhub.watchdog.plist
-->
<plist version="1.0">
<dict>
  <!-- Unique identifier for launchd -->
  <key>Label</key>
  <string>com.productionhub.watchdog</string>

  <!-- The watchdog script path -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__HUB_DIR__/service/watchdog.sh</string>
  </array>

  <!-- Working directory -->
  <key>WorkingDirectory</key>
  <string>__HUB_DIR__</string>

  <!-- Environment variables -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>__NODE_DIR__:/usr/local/bin:/usr/bin:/bin</string>
    <key>HUB_DIR</key>
    <string>__HUB_DIR__</string>
    <key>NODE_BIN</key>
    <string>__NODE_BIN__</string>
    <key>HEALTH_PORT</key>
    <string>8080</string>
  </dict>

  <!-- Start at login -->
  <key>RunAtLoad</key>
  <true/>

  <!-- Keep the watchdog alive (launchd restarts it if it crashes) -->
  <key>KeepAlive</key>
  <true/>

  <!-- Don't restart faster than every 30 seconds -->
  <key>ThrottleInterval</key>
  <integer>30</integer>

  <!-- Logging -->
  <key>StandardOutPath</key>
  <string>__HUB_DIR__/logs/watchdog.log</string>
  <key>StandardErrorPath</key>
  <string>__HUB_DIR__/logs/watchdog.log</string>
</dict>
</plist>
